# mirror-fedora.cron
# rsync from upstream fedora mirror every 5 hours
45 */5 * * *  root /usr/bin/mirror-fedora 2&>1 /dev/null
